# integrations/pcap_analyzers/crontab
# once a day
0 0 * * * suricata-update update-sources 2>&1 >> /var/log/cron.log && suricata-update 2>&1 >> /var/log/cron.log
# todo-ricaricare anche la socket