# integrations/pcap_analyzers/crontab
# remember to keep the newline otherwise the Cron breaks
# once a day
0 0 * * * /update_signatures.sh 2>&1 >> /var/log/cron.log
